# systemd/tunnelgate-proxy.service
[Unit]
Description=TunnelGate Proxy Core
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/tunnelgate start
Restart=always
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
